# service/upstart/ckb-daemon.conf
# ckb-daemon service file
description "Corsair keyboard daemon"
author "ccMSC <qccmsc@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

exec /usr/bin/ckb-daemon
