# debian/ckb-next.service
[Unit]
Description=Corsair Keyboard Daemon

[Service]
ExecStart=/usr/bin/ckb-next-daemon
Restart=always
StandardOutput=syslog
StandardError=inherit

[Install]
WantedBy=multi-user.target
