# linux/systemd/ckb-next-daemon.service.in
[Unit]
Description=Corsair Keyboards and Mice Daemon

[Service]
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/ckb-next-daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
